# node_filesystem/etc/service/bird/run
#!/bin/sh
exec 2>&1
exec bird -R -s /sockets/bird.ctl -d -c /config/bird.cfg